# .config/systemd/user/display-switch.service
[Unit]
Description=USB-triggered display switch

[Service]
ExecStart=%h/.nix-profile/bin/display_switch

[Install]
WantedBy=default.target
