# .config/systemd/user/mbsync.timer
[Unit]
Description=Start mbsync

[Timer]
OnBootSec=10min
OnUnitActiveSec=10min
Unit=mbsync.service

[Install]
WantedBy=default.target
