# .config/systemd/user/test.service
[Unit]
Description=Test service

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'echo $PATH'
